# deploy/supervisord.conf
; =============================================================================
; Parro — Supervisor 配置（用于 Docker 容器内进程管理）
; 同时启动 Flask (gunicorn) + Pipeline Poller
; =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:flask]
command=gunicorn server:app --bind 0.0.0.0:8899 --workers 4 --worker-class sync --timeout 120 --access-logfile - --error-logfile -
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[program:poller]
command=python3 realtime_poller.py
directory=/app/pipeline
autostart=true
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
killasgroup=true
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
